# _.cron
*/15 * * * * sudo /etc/init.d/rumax-wss check